# enable-ir.service
[Unit]
Description=Map IR codes

[Service]
ExecStart=ir-keytable -c -w /etc/rc_keymaps/grunding

[Install]
WantedBy=multi-user.target